# ansible/roles/gps/templates/gps-init.j2
#!/bin/sh /etc/rc.common

START=15
STOP=90

# WM1302 Pi HAT + Quectel L76K(B)
GPIO_RST={{ gps_gpio_rst }}       # GPS_RST (active high pulse)
GPIO_WAKE={{ gps_gpio_wake }}      # GPS_WAKE_UP (high = standby, low = active)
TTY_DEV="{{ gps_tty_device }}"
BAUD={{ gps_baud }}

log() {
    logger -t gps-init "$@"
    [ "$VERBOSE" = "1" ] && echo "gps-init: $@"
}

export_gpio() {
    local gpio="$1"
    [ -d "/sys/class/gpio/gpio${gpio}" ] && return 0
    echo "$gpio" > /sys/class/gpio/export 2>/dev/null || true
}

set_gpio_out() {
    local gpio="$1"
    export_gpio "$gpio"
    echo "out" > "/sys/class/gpio/gpio${gpio}/direction" 2>/dev/null || true
}

set_gpio_val() {
    local gpio="$1"
    local val="$2"
    echo "$val" > "/sys/class/gpio/gpio${gpio}/value" 2>/dev/null || true
}

init_gps_gpio() {
    log "configuring GPS GPIOs (RST=${GPIO_RST}, WAKE=${GPIO_WAKE})"

    # Make sure both pins are outputs
    set_gpio_out "$GPIO_RST"
    set_gpio_out "$GPIO_WAKE"

    # Ensure reset is low initially
    set_gpio_val "$GPIO_RST" 0

    # Wake GPS (WAKE low = active, high = standby)
    set_gpio_val "$GPIO_WAKE" 0

    # Short reset pulse: high for ~20ms, then low
    log "pulsing GPS reset on GPIO${GPIO_RST}"
    set_gpio_val "$GPIO_RST" 1
    sleep 0.02
    set_gpio_val "$GPIO_RST" 0

    # Small settle delay
    sleep 0.1

    log "GPIO states:"
    if [ -r "/sys/class/gpio/gpio${GPIO_RST}/value" ]; then
        log "  GPIO${GPIO_RST} (RST)     = $(cat /sys/class/gpio/gpio${GPIO_RST}/value)"
    fi
    if [ -r "/sys/class/gpio/gpio${GPIO_WAKE}/value" ]; then
        log "  GPIO${GPIO_WAKE} (WAKEUP) = $(cat /sys/class/gpio/gpio${GPIO_WAKE}/value)"
    fi
}

config_uart() {
    if [ ! -e "$TTY_DEV" ]; then
        log "WARNING: $TTY_DEV does not exist, skipping stty"
        return 0
    fi

    log "configuring UART $TTY_DEV at ${BAUD} 8N1 (no flow control)"
    stty -F "$TTY_DEV" "$BAUD" cs8 -cstopb -parenb -ixon 2>/dev/null || log "stty failed on $TTY_DEV"
}

gps_on() {
    set_gpio_out "$GPIO_RST"
    set_gpio_out "$GPIO_WAKE"
    
    # Reset GPS (ensure LOW, pulse HIGH for 20ms, return to LOW)
    set_gpio_val "$GPIO_RST" 0
    sleep 0.01
    set_gpio_val "$GPIO_RST" 1
    sleep 0.02
    set_gpio_val "$GPIO_RST" 0
    
    # Wake GPS (LOW = active)
    set_gpio_val "$GPIO_WAKE" 0
    
    echo "GPS turned ON (active mode)"
    echo "GPIO ${GPIO_RST} (RST): Reset pulse completed, now LOW (0)"
    echo "GPIO ${GPIO_WAKE} (WAKE_UP): Set to LOW (0) - active mode"
}

gps_off() {
    set_gpio_out "$GPIO_RST"
    set_gpio_out "$GPIO_WAKE"
    
    # Ensure GPIO RST is LOW (not in reset state)
    set_gpio_val "$GPIO_RST" 0
    
    # Put GPS in standby (HIGH = standby)
    set_gpio_val "$GPIO_WAKE" 1
    
    echo "GPS turned OFF (standby mode)"
    echo "GPIO ${GPIO_RST} (RST): Set to LOW (0) - normal operation"
    echo "GPIO ${GPIO_WAKE} (WAKE_UP): Set to HIGH (1) - standby mode"
}

gps_status() {
    if [ ! -d "/sys/class/gpio/gpio${GPIO_WAKE}" ]; then
        echo "GPS Status: Not configured (GPIO${GPIO_WAKE} not exported)"
        return 2
    fi
    
    local wake_val=$(cat "/sys/class/gpio/gpio${GPIO_WAKE}/value" 2>/dev/null)
    local rst_val=""
    if [ -d "/sys/class/gpio/gpio${GPIO_RST}" ]; then
        rst_val=$(cat "/sys/class/gpio/gpio${GPIO_RST}/value" 2>/dev/null)
    fi
    
    if [ "$wake_val" = "0" ]; then
        echo "GPS Status: ON (active)"
    elif [ "$wake_val" = "1" ]; then
        echo "GPS Status: OFF (standby)"
    else
        echo "GPS Status: Unknown state"
        return 2
    fi
    
    if [ -n "$rst_val" ]; then
        echo "GPIO ${GPIO_RST} (RST) value: ${rst_val} - 0=normal, 1=reset"
    else
        echo "GPIO ${GPIO_RST} (RST): Not exported"
    fi
    echo "GPIO ${GPIO_WAKE} (WAKE_UP) value: ${wake_val} - 0=active, 1=standby"
    echo ""
    echo "GPIO Control:"
    echo "  GPIO ${GPIO_RST} (RST): LOW (0) = normal operation, HIGH (1) = reset"
    echo "  GPIO ${GPIO_WAKE} (WAKE_UP): LOW (0) = GPS Active (ON), HIGH (1) = GPS Standby (OFF)"
}

start() {
    log "starting GPS init"
    init_gps_gpio
    config_uart
    log "GPS init done; NMEA should be available on $TTY_DEV at ${BAUD} baud"
}

stop() {
    # nothing to tear down; GPS can keep running
    log "stop called (no action)"
}

# Handle custom commands (on/off/state) when script is called directly
if [ -n "$1" ] && [ "$1" != "start" ] && [ "$1" != "stop" ] && [ "$1" != "restart" ] && [ "$1" != "reload" ] && [ "$1" != "enable" ] && [ "$1" != "disable" ] && [ "$1" != "status" ] && [ "$1" != "enabled" ]; then
    VERBOSE=1
    case "$1" in
        on)
            gps_on
            exit $?
            ;;
        off)
            gps_off
            exit $?
            ;;
        state)
            gps_status
            exit $?
            ;;
        *)
            echo "Usage: $0 {start|stop|restart|reload|enable|disable|status|enabled|on|off|state}"
            echo ""
            echo "Init script commands:"
            echo "  start   - Initialize GPS at boot"
            echo "  stop    - Stop GPS (no action)"
            echo ""
            echo "Runtime control:"
            echo "  on      - Turn GPS ON (active mode)"
            echo "  off     - Turn GPS OFF (standby mode)"
            echo "  state   - Check GPS status"
            exit 1
            ;;
    esac
fi
